# build/Linux/opt/devatserv/share/sys/remove_gui_app.service
[Unit]
Description=Remove DevAtServ'GUI App Service

[Service]
Type=exec
ExecStart=/opt/devatserv/share/util/remove_gui_app.sh

[Install]
WantedBy=multi-user.target