# build/Linux/opt/share/sys/remove_gui_app.service
[Unit]
Description=Remove DevAtServ'GUI App Service

[Service]
Type=exec
ExecStart=/opt/share/util/remove_gui_app.sh
